# files/common/etc/cron.weekly/cleanscratch
#!/bin/sh

df -h /scratch
find /scratch/ -xdev -atime +1460 -delete
df -h /scratch
