# files/common/etc/cron.weekly/cleanunowned
#!/bin/sh
df -h /scratch
sleep 60
id testingn00 && find /scratch/ -xdev -nouser -delete
df -h /scratch
